# deploy/vantage-dataengine.timer
# Daily timer for the Vantage data-engine refresh (§9).
#
# Install:
#   sudo cp deploy/vantage-dataengine.service /etc/systemd/system/
#   sudo cp deploy/vantage-dataengine.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vantage-dataengine.timer
#
# Check / run / inspect:
#   systemctl list-timers vantage-dataengine     # when it next fires
#   sudo systemctl start vantage-dataengine      # run it now, once
#   journalctl -u vantage-dataengine -n 50       # what it did
#
# DAILY is safe because the freshness gate (--min-age-days 30) makes runs no-ops
# unless an area is new or stale — a day with nothing due costs nothing. What
# daily buys over weekly: a camera added today has its area context tomorrow
# morning, not next Sunday. Spend stays bounded by the gate + per-run caps.

[Unit]
Description=Daily Vantage data-engine refresh
Requires=vantage-dataengine.service

[Timer]
OnCalendar=03:30
# If the box was off when it was due, run shortly after boot rather than skipping.
Persistent=true
# Spread the load / avoid a thundering herd against Apify.
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
