# deploy/vantage.service
# systemd unit for the Vantage FastAPI app (uvicorn), co-located on the box.
# Install:
#   sudo cp deploy/vantage.service /etc/systemd/system/vantage.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vantage
#   systemctl status vantage        # check it's running
#   journalctl -u vantage -f        # follow logs (first run downloads models)
#
# Binds 127.0.0.1:8000 (via production.env) so only Caddy reaches it.

[Unit]
Description=Vantage API (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/vantage
EnvironmentFile=-/home/ubuntu/vantage/alibi/config/production.env
# Bind 127.0.0.1:8600 explicitly (Caddy fronts it). Port 8600 avoids the
# existing services on this box — AIKit is on :8000, the Node app on :3001,
# the Nodes on :4101-4105. The app's own default is :8000, so we override it
# here rather than relying on config/env.
ExecStart=/home/ubuntu/vantage/.venv/bin/python -m uvicorn alibi.alibi_api:app --host 127.0.0.1 --port 8600
Restart=always
RestartSec=5
# ML models (D-FINE, OSNet, FastALPR) load into RAM on first request; allow time.
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
